# solrcloud5/supervisord.conf
[supervisord]
nodaemon=true

[program:solr]
command=/opt/bin/start-server.sh

